# custom/splash-cycler/custom-splash-screens.timer
[Unit]
Description=Timer to set custom splash screens

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=multi-user.target

