# custom/splash-cycler/custom-splash-screens.timer
[Unit]
Description=Timer to set custom splash screens

[Timer]
OnCalendar=TIMERSPECIFICATION

[Install]
WantedBy=multi-user.target

